# learning_log-supervisor.conf
[program:learning_log]   ##修改成你的项目名称
command=/bin/learning_log-gunicorn.sh  
environment=DJANGO_ENV="prod"
#user=wang ##修改成你的用户
autostart=true
autorestart=true
redirect_stderr=True
